# backend/ec2/bric-a-brac-drf/my_services/bricabrac.service
[Unit]
Description=bricabrac service

[Service]
User=root
Group=root
WorkingDirectory=/
EnvironmentFile=/bricabrac.env
ExecStart=/usr/local/bin/gunicorn bricabrac.wsgi:application --config /gunicorn.conf.py
